# proxy-server/linux/demergi.service
[Unit]
Description=Demergi daemon
After=network.target

[Service]
User=root
ExecStart=demergi

[Install]
WantedBy=multi-user.target
